# root/home/victor/.config/systemd/user/checkpoint.timer
[Timer]
OnUnitActiveSec=25m
AccuracySec=1s
Unit=checkpoint.service

[Install]
WantedBy=timers.target

