# root/home/victor/.checkpoint/checkpoint.service
# ~/.config/systemd/user/checkpoint.service
[Service]
Type=oneshot
ExecStart=%h/.checkpoint/prompt.sh

